# ops/systemd/halcyon-price-relay.service
[Unit]
Description=Halcyon Pyth price relay (Hermes → devnet PriceUpdateV2 accounts)
Documentation=file:///etc/halcyon/config/price_relay.json
Documentation=https://github.com/DJB8787/colosseumfinal/blob/main/keepers/price_relay/README.md
After=network-online.target
Wants=network-online.target
PartOf=halcyon-keepers.target

[Service]
Type=simple
User=halcyon
Group=halcyon
WorkingDirectory=/opt/halcyon/price_relay
EnvironmentFile=-/etc/halcyon/env/price_relay.env
ExecStart=/usr/bin/env node --import tsx/esm src/index.ts --config /etc/halcyon/config/price_relay.json

Restart=always
RestartSec=5s
StartLimitBurst=5
StartLimitIntervalSec=300

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/halcyon
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

StandardOutput=journal
StandardError=journal
SyslogIdentifier=halcyon-price-relay

[Install]
WantedBy=halcyon-keepers.target
