# ops/hetzner-deploy/systemd/halcyon-flagship-sigma-keeper.timer
[Unit]
Description=Run the Halcyon flagship sigma keeper daily

[Timer]
Unit=halcyon-flagship-sigma-keeper.service
OnCalendar=*-*-* 00:10:00
Persistent=true
AccuracySec=1s

[Install]
WantedBy=timers.target
