# ops/hetzner-deploy/systemd/halcyon-fire-reduced-ops.timer
[Unit]
Description=Run the Halcyon SOL reduced-operator uploader every minute

[Timer]
Unit=halcyon-fire-reduced-ops.service
OnCalendar=*-*-* *:*:40
Persistent=true
AccuracySec=1s

[Install]
WantedBy=timers.target
