# ops/hetzner-deploy/systemd/halcyon-delta-keeper.service
[Unit]
Description=Halcyon flagship aggregate delta keeper
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/halcyon
EnvironmentFile=/etc/halcyon/env
ExecStart=/opt/halcyon/bin/delta_keeper --config /etc/halcyon/config/delta_keeper.json
Restart=always
RestartSec=10s
StartLimitIntervalSec=300
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=halcyon-delta-keeper

[Install]
WantedBy=multi-user.target
