# ops/hetzner-deploy/systemd/halcyon-update-ewma-il.service
[Unit]
Description=Halcyon IL EWMA updater
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/halcyon
EnvironmentFile=/etc/halcyon/env
ExecStart=/usr/bin/env bash -lc 'exec /opt/halcyon/bin/halcyon --rpc "$HELIUS_DEVNET_RPC" --keypair "$IL_EWMA_KEYPAIR" update-ewma --product il --oracle-price "$PYTH_SOL_ACCOUNT"'
StandardOutput=journal
StandardError=journal
SyslogIdentifier=halcyon-update-ewma-il
